# ipservice.service
Description=IP Service

[Service]
ExecStart=/usr/local/bin/ip-service -listen :8000 -log /dev/null

[Install]
WantedBy=multi-user.target
